# wunderground/wu-update.cron
# $Id: wu-update.cron,v 1.1 2003/06/02 19:59:25 wdp Exp $
#
# Example crontab entry for uploading data to wunderground every 5 mins
#

*/5 * * * * $HOME/weather/do-wu-update.sh

